# images/cron/crontab
*/15 * * * * chmod 775 -R /data/nextcloud/data
0 * * * * docker exec nextcloud occ files:scan --all
