# images/rsnapshot/crontab
0 */4 * * *	/usr/bin/rsnapshot hourly
30 23 * * *	/usr/bin/rsnapshot daily
