# images/samba/supervisord.conf
[supervisord]
nodaemon=true

[program:samba]
command=/usr/sbin/smbd --daemon --foreground --log-stdout

[program:nmbd]
command=/usr/sbin/nmbd --daemon --foreground --log-stdout
